# deploy/macu-ttyd/macu-ttyd.service
# systemd --user unit for the MACU terminal — ttyd serving a tmux session running
# Claude Code. It backs Studio's right-hand TERMINAL drawer (the iframe at :7682).
# The setup-macu-channel skill resolves absolute ttyd/claude paths + a PATH,
# substitutes the __PLACEHOLDERS__, installs this to ~/.config/systemd/user/, then
# `systemctl --user enable --now macu-ttyd`. (On WSL / no user-systemd the skill
# falls back to launching ttyd under nohup.)
#
# SECURITY: ttyd serves an INTERACTIVE Claude session with NO AUTH — that's shell
# access (via Claude's tools) to anyone who can reach the port. It binds LOOPBACK
# (127.0.0.1) ONLY by default (the `-i 127.0.0.1` below), so it's reachable from the
# machine running Studio — including WSL from the Windows host, via localhost
# forwarding — but NOT from the LAN. To use the drawer from ANOTHER machine over the
# LAN, remove `-i 127.0.0.1` — but you MUST also add `-c user:pass` (basic auth) in
# the same edit, or you are publishing an unauthenticated root shell to your whole
# network. Then point the frontend at the LAN host with VITE_TERMINAL_URL. NEVER
# expose :7682 on the WAN.
[Unit]
Description=MACU terminal (ttyd + tmux running Claude Code, for Studio's TERMINAL drawer)
After=default.target

[Service]
Type=simple
Environment=PATH=__PATH__
ExecStart=__TTYD_BIN__ -i 127.0.0.1 -W -p 7682 tmux new-session -A -s claude __CLAUDE_BIN__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
