# deploy/macu-render.service
[Unit]
# TEMPLATE — __PLACEHOLDERS__ are substituted for this machine by
# deploy/install-systemd.sh (repo path, run user). Don't cp it verbatim.
Description=MACU Render HTTP Service (drives pipeline/run.py)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__REPO__/pipeline
# Optional repo-root .env (leading '-' = no-op if absent).
EnvironmentFile=-__REPO__/.env
# Run under the Studio venv so serve.py → run.py have the pipeline deps
# (Pillow, python-dotenv) and a consistent interpreter.
ExecStart=__REPO__/studio/.venv/bin/python __REPO__/pipeline/serve.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/macu-render.log
StandardError=append:/var/log/macu-render.log

[Install]
WantedBy=multi-user.target
