# deploy/macu-render.service
[Unit]
Description=MACU Render HTTP Service (drives pipeline/run.py)
Documentation=https://10.0.0.245:8002/projects/3 (SSA-90)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=mayorawesome
Group=mayorawesome
WorkingDirectory=/mnt/storage/shares/MACU/pipeline
ExecStart=/usr/bin/python3 /mnt/storage/shares/MACU/pipeline/serve.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/macu-render.log
StandardError=append:/var/log/macu-render.log

[Install]
WantedBy=multi-user.target
