# studio/systemd/macu-studio.service
[Unit]
# TEMPLATE — __PLACEHOLDERS__ are substituted for this machine by
# deploy/install-systemd.sh (repo path, shares dir, run user). Don't cp it verbatim.
Description=MACU Studio (FastAPI + SPA on :8774, drives macu-render)
After=network-online.target macu-render.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__REPO__/studio
Environment=MACU_RENDER_URL=http://127.0.0.1:8773
Environment=MACU_EPISODES=__SHARES__/episodes
Environment=PYTHONUNBUFFERED=1
# Studio has no auth — bind loopback by default. Set MACU_STUDIO_HOST=0.0.0.0 in
# .env to reach it from another machine on a trusted LAN.
Environment=MACU_STUDIO_HOST=127.0.0.1
Environment=MACU_STUDIO_PORT=8774
# Optional repo-root .env (the leading '-' makes it a no-op if absent). It is read
# AFTER the inline Environment= lines above, so values it sets (e.g. MACU_STUDIO_HOST)
# override the defaults here.
EnvironmentFile=-__REPO__/.env
ExecStart=__REPO__/studio/.venv/bin/uvicorn macu_studio.main:app --host ${MACU_STUDIO_HOST} --port ${MACU_STUDIO_PORT}
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/macu-studio.log
StandardError=append:/var/log/macu-studio.log

[Install]
WantedBy=multi-user.target
